# files/etc/init.d/mac-randomizer
#!/bin/sh /etc/rc.common
# 设备指纹完全随机化服务 v3.2
# 用途：让 Temu/TikTok 等平台识别为全新设备
# 触发：每次重启自动执行
# 作者：SimplexY-tech
# 
# 实现方式：通过启动脚本在系统启动时修改配置文件（方案1）
# 稳定性：优先（不依赖编译时注入，使用标准 UCI 配置）

START=11
STOP=89

LOG_TAG="device-fingerprint"

# ==============================
# 日志函数
# ==============================

log_info() {
    logger -t "$LOG_TAG" -s "INFO: $1"
}

log_error() {
    logger -t "$LOG_TAG" -s "ERROR: $1"
}

# ==============================
# 1. 生成随机前缀（2字母）
# ==============================

generate_random_prefix() {
    # 返回 2 个大写随机字母
    # 用于 Hostname 和 5G SSID 的前缀
    local letters="ABCDEFGHIJKLMNOPQRSTUVWXYZ"
    local letter1=${letters:$((RANDOM % 26)):1}
    local letter2=${letters:$((RANDOM % 26)):1}
    echo "${letter1}${letter2}"
}

# ==============================
# 2. 真实厂商 MAC 前缀库
# ==============================

get_random_vendor_mac_prefix() {
    # 真实路由器厂商 OUI 前缀库
    local vendors=(
        # TP-Link
        "50:3E:AA" "60:E3:27" "F4:F2:6D" "30:B4:9E" "14:CC:20" "C0:25:E9"
        # Linksys (Cisco)
        "00:18:F8" "00:1D:7E" "00:23:69" "08:62:66" "10:FE:ED"
        # ASUS
        "60:A4:4C" "AC:22:0B" "E0:3F:49" "84:A9:3E" "AC:9E:17"
        # D-Link
        "00:05:5D" "00:13:46" "14:CF:92" "20:AB:37" "AC:F1:DF"
        # Netgear
        "00:0F:B5" "00:1B:2F" "20:4E:7F" "A0:40:A0" "44:94:FC"
        # Huawei
        "00:11:6B" "00:25:68" "68:A8:6D" "18:82:8C" "F8:DB:88"
        # Xiaomi
        "08:56:5B" "34:CE:00" "5C:6A:8F"
        # Tenda
        "F8:D1:11" "A0:2B:B8"
    )
    
    local idx=$(( RANDOM % ${#vendors[@]} ))
    echo "${vendors[$idx]}"
}

generate_real_mac() {
    # 生成符合真实厂商的 MAC 地址
    local vendor_prefix=$(get_random_vendor_mac_prefix)
    local device_suffix=$(printf '%02x:%02x:%02x' \
        $(( RANDOM & 0xFF )) \
        $(( RANDOM & 0xFF )) \
        $(( RANDOM & 0xFF )))
    
    echo "$vendor_prefix:$device_suffix"
}

# ==============================
# 3. 有线网络 MAC 随机化
# ==============================

randomize_wired_macs() {
    log_info "===== 有线网络 MAC 随机化（真实厂商） ====="
    
    # WAN 口 (eth0)
    if ip link show eth0 >/dev/null 2>&1; then
        ip link set eth0 down 2>/dev/null
        sleep 0.5
        NEW_MAC=$(generate_real_mac)
        if ! ip link set eth0 address "$NEW_MAC" 2>/dev/null; then
            log_error "Failed to set eth0 MAC: $NEW_MAC"
            return 1
        fi
        ip link set eth0 up 2>/dev/null
        log_info "eth0 MAC → $NEW_MAC ✓"
    else
        log_error "eth0 not found"
    fi
    
    # LAN 口 (eth1)
    if ip link show eth1 >/dev/null 2>&1; then
        ip link set eth1 down 2>/dev/null
        sleep 0.5
        NEW_MAC=$(generate_real_mac)
        if ! ip link set eth1 address "$NEW_MAC" 2>/dev/null; then
            log_error "Failed to set eth1 MAC: $NEW_MAC"
            return 1
        fi
        ip link set eth1 up 2>/dev/null
        log_info "eth1 MAC → $NEW_MAC ✓"
    else
        log_error "eth1 not found"
    fi
    
    # 网桥 (br-lan)
    if ip link show br-lan >/dev/null 2>&1; then
        NEW_MAC=$(generate_real_mac)
        if ! ip link set br-lan address "$NEW_MAC" 2>/dev/null; then
            log_error "Failed to set br-lan MAC: $NEW_MAC"
            return 1
        fi
        log_info "br-lan MAC → $NEW_MAC ✓"
    else
        log_error "br-lan not found"
    fi
    
    log_info "有线网络 MAC 随机化完成"
}

# ==============================
# 4. WiFi 配置随机化
# ==============================

randomize_wifi_config() {
    log_info "===== WiFi 配置随机化 ====="
    
    local uci_changed=0
    
    # 获取前缀（Hostname 和 5G SSID 前缀相同）
    RANDOM_PREFIX=$(generate_random_prefix)
    
    # 2.4G WiFi: MAC随机 + SSID固定 + 密码固定
    if uci get wireless.@wifi-device[0] >/dev/null 2>&1; then
        MAC_2G=$(generate_real_mac)
        
        uci set wireless.@wifi-device[0].macaddr="$MAC_2G"
        uci set wireless.@wifi-iface[0].ssid="OWRT"
        uci set wireless.@wifi-iface[0].key="12345678"
        uci set wireless.@wifi-iface[0].encryption="psk2+ccmp"
        
        log_info "2.4G: SSID=OWRT | MAC=$MAC_2G ✓"
        uci_changed=1
    fi
    
    # 5G WiFi: MAC随机 + SSID随机(规则) + 密码固定
    # SSID 前缀与 Hostname 相同
    if uci get wireless.@wifi-device[1] >/dev/null 2>&1; then
        MAC_5G=$(generate_real_mac)
        # 5G SSID 格式：前缀(AB) + 月日(0531) + 时分(1430)
        DATE_PART=$(date +%m%d)
        TIME_PART=$(date +%H%M)
        SSID_5G="${RANDOM_PREFIX}${DATE_PART}${TIME_PART}"
        
        uci set wireless.@wifi-device[1].macaddr="$MAC_5G"
        uci set wireless.@wifi-iface[1].ssid="$SSID_5G"
        uci set wireless.@wifi-iface[1].key="12345678"
        uci set wireless.@wifi-iface[1].encryption="psk2+ccmp"
        
        log_info "5G: SSID=$SSID_5G | MAC=$MAC_5G ✓"
        uci_changed=1
    fi
    
    if [ $uci_changed -eq 1 ]; then
        uci commit wireless
        sleep 1
        wifi down && sleep 2 && wifi up
        sleep 3
        log_info "WiFi 配置已提交并重启 ✓"
    else
        log_error "未找到WiFi设备配置"
        return 1
    fi
}

# ==============================
# 5. 系统标识随机化
# ==============================

randomize_system_ids() {
    log_info "===== 系统标识随机化 ====="
    
    # Hostname: 使用相同的前缀 + 月日
    # 格式: AB0531 (2字母前缀 + 月日)
    RANDOM_PREFIX=$(generate_random_prefix)
    DATE_PART=$(date +%m%d)
    RANDOM_HOSTNAME="${RANDOM_PREFIX}${DATE_PART}"
    
    uci set system.@system[0].hostname="$RANDOM_HOSTNAME"
    uci commit system
    hostname "$RANDOM_HOSTNAME"
    
    log_info "Hostname → $RANDOM_HOSTNAME (前缀与5G SSID相同) ✓"
    
    # 系统描述
    uci set system.@system[0].description="OpenWrt"
    uci commit system
    
    # 清理设备 ID 文件
    [ -f /etc/device.id ] && rm -f /etc/device.id
    
    log_info "系统标识随机化完成 ✓"
}

# ==============================
# 6. UPnP 设备型号隐藏
# ==============================

hide_upnp_device_info() {
    log_info "===== UPnP 设备信息隐藏 ====="
    
    if uci get upnp >/dev/null 2>&1; then
        uci set upnp.config.enabled=0
        uci commit upnp
        /etc/init.d/upnpd stop >/dev/null 2>&1
        log_info "UPnP 已禁用 ✓"
    fi
    
    # SSDP 广播隐藏
    echo "1" > /proc/sys/net/ipv4/icmp_echo_ignore_broadcasts 2>/dev/null
    log_info "SSDP 广播已隐藏 ✓"
}

# ==============================
# 7. HTTP Server 响应头隐藏
# ==============================

hide_http_server_banner() {
    log_info "===== HTTP Server 指纹隐藏 ====="
    
    local uhttpd_conf="/etc/config/uhttpd"
    
    if [ -f "$uhttpd_conf" ]; then
        if ! uci get uhttpd.main.http_keepalive >/dev/null 2>&1; then
            uci set uhttpd.main.http_keepalive=1
        fi
        
        uci commit uhttpd
        /etc/init.d/uhttpd restart >/dev/null 2>&1
        log_info "HTTP Server 响应头已隐藏 ✓"
    fi
}

# ==============================
# 8. SSH/Telnet 横幅隐藏
# ==============================

hide_ssh_telnet_banner() {
    log_info "===== SSH/Telnet 横幅隐藏 ====="
    
    # SSH 横幅
    local ssh_banner="/etc/banner.net"
    if [ -f "$ssh_banner" ]; then
        > "$ssh_banner"
        echo "Welcome" > "$ssh_banner"
        log_info "SSH 横幅已修改 ✓"
    fi
    
    # SSH 隐藏版本
    local ssh_config="/etc/ssh/sshd_config"
    if [ -f "$ssh_config" ]; then
        if ! grep -q "VersionAddendum" "$ssh_config"; then
            echo "VersionAddendum none" >> "$ssh_config"
        fi
        /etc/init.d/sshd restart >/dev/null 2>&1
        log_info "SSH 版本信息已隐藏 ✓"
    fi
    
    # Telnet 横幅
    local telnet_banner="/etc/banner"
    if [ -f "$telnet_banner" ]; then
        > "$telnet_banner"
        echo "Welcome" > "$telnet_banner"
        log_info "Telnet 横幅已修改 ✓"
    fi
}

# ==============================
# 9. 网络指纹清理
# ==============================

clean_network_fingerprints() {
    log_info "===== 网络指纹清理 ====="
    
    # DHCP 隐匿配置
    if uci get dhcp.@dnsmasq[0] >/dev/null 2>&1; then
        uci set dhcp.@dnsmasq[0].rebind_protection=0
        uci set dhcp.@dnsmasq[0].localise_queries=0
        uci set dhcp.@dnsmasq[0].nonegative_ttl=1
        uci set dhcp.@dnsmasq[0].log_queries=0
        uci commit dhcp
        /etc/init.d/dnsmasq restart >/dev/null 2>&1
        log_info "DHCP 隐匿配置已应用 ✓"
    fi
    
    # mDNS 禁用
    if command -v avahi-daemon >/dev/null 2>&1; then
        /etc/init.d/avahi-daemon stop >/dev/null 2>&1
        log_info "mDNS 已禁用 ✓"
    fi
}

# ==============================
# 10. 存储指纹清理
# ==============================

clean_storage_fingerprints() {
    log_info "===== 存储指纹清理 ====="
    
    # 清理系统日志
    rm -f /var/log/messages /var/log/syslog /var/log/*.log* 2>/dev/null
    > /var/log/messages
    log_info "系统日志已清理 ✓"
    
    # 清理临时文件
    rm -rf /tmp/* 2>/dev/null
    rm -rf /var/tmp/* 2>/dev/null
    log_info "临时文件已清理 ✓"
    
    # 清理 DHCP/DNS 日志
    [ -f /var/log/dhcp.log ] && > /var/log/dhcp.log
    [ -f /var/log/dnsmasq ] && > /var/log/dnsmasq
    
    # 清理 Luci 会话
    rm -rf /var/luci-sessions/* 2>/dev/null
    
    # 内存缓存清理
    sync
    echo 3 > /proc/sys/vm/drop_caches 2>/dev/null
    
    log_info "存储指纹清理完成 ✓"
}

# ==============================
# 11. 网络参数随机化
# ==============================

randomize_network_params() {
    log_info "===== 网络参数随机化 ====="
    
    # MTU 随机化
    RANDOM_MTU=$(( 1492 + (RANDOM % 17) ))
    ip link set dev eth0 mtu $RANDOM_MTU 2>/dev/null || true
    log_info "MTU 值 → $RANDOM_MTU ✓"
    
    # TCP 时间戳
    echo 1 > /proc/sys/net/ipv4/tcp_timestamps 2>/dev/null
    
    # IP 分片
    echo 1 > /proc/sys/net/ipv4/ip_no_pmtu_disc 2>/dev/null
    
    log_info "网络参数随机化完成 ✓"
}

# ==============================
# 12. 主程序入口
# ==============================

start() {
    log_info "╔═════════════════════════════════════════════════╗"
    log_info "║ 设备指纹完全随机化服务 v3.2 启动                  ║"
    log_info "║ 实现方式：启动脚本修改 UCI 配置（方案1-最稳定）   ║"
    log_info "║ 触发时机：系统启动时自动执行                     ║"
    log_info "╚═════════════════════════════════════════════════╝"
    
    # 执行所有随机化模块
    log_info ""
    
    randomize_wired_macs || log_error "有线 MAC 随机化失败"
    log_info ""
    
    randomize_wifi_config || log_error "WiFi 配置失败"
    log_info ""
    
    randomize_system_ids || log_error "系统标识随机化失败"
    log_info ""
    
    hide_upnp_device_info || log_error "UPnP 隐藏失败"
    hide_http_server_banner || log_error "HTTP Server 隐藏失败"
    hide_ssh_telnet_banner || log_error "SSH/Telnet 隐藏失败"
    log_info ""
    
    clean_network_fingerprints || log_error "网络指纹清理失败"
    clean_storage_fingerprints || log_error "存储指纹清理失败"
    log_info ""
    
    randomize_network_params || log_error "网络参数随机化失败"
    log_info ""
    
    # 获取最终状态
    FINAL_HOSTNAME=$(hostname)
    FINAL_5G_SSID=$(uci get wireless.@wifi-iface[1].ssid 2>/dev/null)
    FINAL_2G_MAC=$(ip link show | grep -A1 "wlan0" | grep "link/ether" | awk '{print $2}' 2>/dev/null)
    
    log_info "╔═════════════════════════════════════════════════╗"
    log_info "║ 设备指纹随机化完成 ✓                              ║"
    log_info "║ 平台将识别为全新设备！                            ║"
    log_info "║                                                 ║"
    log_info "║ Hostname: $FINAL_HOSTNAME                      ║"
    log_info "║ 5G SSID:  $FINAL_5G_SSID                    ║"
    log_info "║                                                 ║"
    log_info "║ 注：前缀(${FINAL_HOSTNAME:0:2})相同，规则化设计   ║"
    log_info "╚═════════════════════════════════════════════════╝"
}

stop() {
    log_info "设备指纹服务已停止"
}
